# deploy-tools/cron.daily
SHELL=/bin/bash
0 0 * * * ~/student-dashboard/dailys.py 2>> ~/student-dashboard/daily_error.log